# supervisord.prod.conf
[supervisord]
nodaemon=true

[program:server]
command=pm2 start server.js -i $NODE_INSTANCES --max-memory-restart 400M
retries = 10
