# .prod/supervisord.prod.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
retries=5

[program:site]
command=node site.js
retries=5

[program:services]
command=node app/services.js
retries=5
